# local-tunnel.service
[Unit]
Description=SSH Tunnels to Droplet
After=network.target

[Service]
EnvironmentFile=/Users/aasim/.config/ssh-tunnel.env
ExecStart=/usr/bin/autossh -M 0 -N -q \
  -o "ServerAliveInterval 60" \
  -o "ServerAliveCountMax 3" \
  -L 3000:localhost:3000 \
  -L 8000:localhost:8000 \
  root@${DROPLET_IP} -i ${SSH_KEY_PATH}
Restart=always
RestartSec=10

[Install]
WantedBy=default.target 